# services/olcrtc.service
[Unit]
Description=OlcRTC SFU Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/olcrtc
ExecStartPre=/usr/bin/mkdir -p /opt/olcrtc/data
ExecStart=/usr/local/bin/olcrtc -mode sfu -carrier wbstream -dns 1.1.1.1:53 -room ROOM_ID -client CLIENT_ID -key KEY -data /opt/olcrtc/data -port PORT
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
